# conf/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=php-fpm
stdout_logfile=/var/log/supervisor/php-fpm-stdout.log
stderr_logfile=/var/log/supervisor/php-fpm-stderr.log
autostart=true
autorestart=true