# conf/supervisord/supervisor.ini

[supervisord]
nodaemon=true

[program:nginx]
command = /usr/sbin/nginx
autostart = true

[program:php-fpm]
command = /usr/bin/php-fpm -F
autostart = true