# conf/supervisord/supervisor.conf

[supervisord]
nodaemon=true

[program:nginx]
command = nginx
autostart = true

[program:php-fpm]
command = php-fpm -F
autostart = true

[program:sshd]
command = /usr/sbin/sshd
autostart = true
